# services.d/tunnelto_server/run
#!/usr/bin/execlineb -S0

multisubstitute
{
    importas -D "" TUNNELTO_TUNNEL_PORT TUNNELTO_TUNNEL_PORT
    importas -D "" TUNNELTO_ALLOWED_HOSTS TUNNELTO_ALLOWED_HOSTS
    importas -D "" TUNNELTO_AWS_ACCESS_KEY_ID TUNNELTO_AWS_ACCESS_KEY_ID
    importas -D "" TUNNELTO_AWS_SECRET_ACCESS_KEY TUNNELTO_AWS_SECRET_ACCESS_KEY
}
emptyenv -p
export RUST_LOG tunnelto_server=debug
export RUST_BACKTRACE 1
export PORT ${TUNNELTO_TUNNEL_PORT}
export ALLOWED_HOSTS ${TUNNELTO_ALLOWED_HOSTS}
export AWS_ACCESS_KEY_ID ${TUNNELTO_AWS_ACCESS_KEY_ID}
export AWS_SECRET_ACCESS_KEY ${TUNNELTO_AWS_SECRET_ACCESS_KEY}
/app/tunnelto_server
